# scripts/launchd/com.tradebot.ngrok.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradebot.ngrok</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ngrok</string>
        <string>http</string>
        <string>5000</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/usmc_reeper_1775/tradebot/logs/ngrok_service.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/usmc_reeper_1775/tradebot/logs/ngrok_service.err.log</string>
</dict>
</plist>
